# etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm --allow-to-run-as-root